# docs/euosint-zone-briefings-sync.timer
[Unit]
Description=Run EUOSINT zone briefings cache sync daily

[Timer]
OnBootSec=3m
OnUnitActiveSec=24h
AccuracySec=1m
Unit=euosint-zone-briefings-sync.service

[Install]
WantedBy=timers.target
